# scripts/renew-cert.cron
0 3 * * * root /usr/local/bin/renew-cert.sh
